# eski sürümler/v2 correct/ai gez/archive_and_backups/uploud/deploy/systemd/gcs-cookie-pull.service
[Unit]
Description=Pull latest cookie payload from GCS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ai-gez
EnvironmentFile=-/opt/ai-gez/.env
Environment=PYTHONUNBUFFERED=1
UMask=0077
NoNewPrivileges=true
ExecStartPre=/opt/ai-gez/.venv/bin/python /opt/ai-gez/server_preflight.py --service gcs-cookie-pull
ExecStart=/opt/ai-gez/.venv/bin/python /opt/ai-gez/redacted_runner.py -- /opt/ai-gez/.venv/bin/python /opt/ai-gez/gcs_cookie_bridge.py pull
StandardOutput=append:/var/log/ai-gez-gcs-pull.log
StandardError=append:/var/log/ai-gez-gcs-pull.log

[Install]
WantedBy=multi-user.target
