# eski sürümler/v2 correct/ai gez/archive_and_backups/deploy_package/deploy/systemd/cookie-refresh.timer
[Unit]
Description=Run Sahibinden cookie refresh check every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
Unit=cookie-refresh.service

[Install]
WantedBy=timers.target
