# eski sürümler/v2 correct/ai gez/archive_and_backups/deploy_package/deploy/systemd/cookie-refresh.service
[Unit]
Description=Sahibinden cookie refresh check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ai-gez
EnvironmentFile=-/opt/ai-gez/.env
ExecStart=/opt/ai-gez/.venv/bin/python /opt/ai-gez/cookie_auto_refresh.py --once --headless --force-refresh-after-min 90 --send-gcs-on-refresh
StandardOutput=append:/var/log/ai-gez-cookie.log
StandardError=append:/var/log/ai-gez-cookie.log

[Install]
WantedBy=multi-user.target
